# roles/duckdns/files/duckdns.service
[Unit]
Description=DuckDNS update job
After=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/duckdns

[Install]
WantedBy=multi-user.target